# scripts/setup/systemd/nest-auto-done.timer
# scripts/setup/systemd/nest-auto-done.timer
# Enable with: sudo systemctl enable --now nest-auto-done.timer
[Unit]
Description=Run the Nest auto-Done gate periodically

[Timer]
OnBootSec=15min
OnUnitActiveSec=30min
Persistent=false

[Install]
WantedBy=timers.target
